# startup/relay-control.conf
description "Relay Control"
author "Bertrand Roussel"

start on runlevel [2345]

exec /root/RelayControl/main.py

